# sys-apps/speakersafetyd/files/speakersafetyd-noroot
#!/sbin/openrc-run

pidfile="/run/speakersafetyd.pid"
command="/usr/bin/speakersafetyd"
command_args="-c /usr/share/speakersafetyd/ -b /var/lib/speakersafetyd/blackbox -m 7"
command_user="speakersafetyd"
supervisor="supervise-daemon"
	
